# fhirhub.service
[Unit]
Description=FHIRHub - Convertisseur HL7 vers FHIR
After=network.target

[Service]
Type=simple
User=fhirhub
WorkingDirectory=/opt/fhirhub
ExecStart=/opt/fhirhub/start.sh
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fhirhub
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target